# supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:odoo]
command=/home/deploy/production/bin/start_openerp -d production
user=deploy
environment=ODOO_CONNECTOR_CHANNELS="root:10"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
